# systemd/ai-sre-agent.service
[Unit]
Description=AI SRE Agent - Automated Server Monitoring
Documentation=https://github.com/melvio/ai-sre-agent
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=melvin
Group=melvin
WorkingDirectory=/home/melvin/server/ai-sre-agent
ExecStart=/usr/bin/python3 /home/melvin/server/ai-sre-agent/agent.py daemon
Restart=always
RestartSec=30

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/home/melvin/.nvm/versions/node/v20.19.6/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
EnvironmentFile=/home/melvin/server/ai-sre-agent/.env

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Logging to systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-sre-agent

[Install]
WantedBy=multi-user.target
