# systemd/signal-chat.service
[Unit]
Description=Signal Chatbot for Home Server
Documentation=https://github.com/melvio/ai-sre-agent
After=network.target

[Service]
Type=simple
User=melvin
Group=melvin
WorkingDirectory=/home/melvin/server/ai-sre-agent
ExecStart=/usr/bin/python3 /home/melvin/server/ai-sre-agent/signal_chat.py
Restart=always
RestartSec=10

# Environment
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/melvin/server/ai-sre-agent/.env

# Logging to systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=signal-chat

[Install]
WantedBy=multi-user.target
